# packaging/systemd/kowalski-core.service
[Unit]
Description=Kowalski OS agent core daemon
# This unit ships in the kowalski-core .deb, which installs the launcher at
# /usr/bin/kow (a symlink into the bundled venv).
ConditionPathExists=/usr/bin/kow

[Service]
Type=simple
ExecStart=/usr/bin/kow serve
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
