# packaging/systemd/kowalski-core.service
[Unit]
Description=Kowalski OS agent core daemon
ConditionPathExists=%h/.local/bin/kow

[Service]
Type=simple
ExecStart=%h/.local/bin/kow serve
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
